# destruct/src/etc/xinetd.d/ctf
service ctf
{
    socket_type = stream
    protocol = tcp
    port = 1337
    wait = no
    user = app
    instances = UNLIMITED
    server = /home/app/wrapper
}
